# backend/webmail-api.service
[Unit]
Description=BNIX Webmail (FastAPI)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=bnix-webmail
Group=bnix-webmail
WorkingDirectory=/opt/bnix-webmail/src/backend
EnvironmentFile=-/etc/bnix-webmail.env
Environment=DATA_DIR=/opt/bnix-webmail/data
Environment=HOST=127.0.0.1
Environment=PORT=8000
ExecStart=/opt/bnix-webmail/venv/bin/python /opt/bnix-webmail/src/backend/main.py
Restart=always
RestartSec=5
PrivateTmp=true
ProtectSystem=full
ReadWritePaths=/opt/bnix-webmail/data /etc/caddy

[Install]
WantedBy=multi-user.target
